# neon_logs/templates/default/flume-ng-service.conf.erb
description "A flume agent"
author  "Mark Desnoyer <desnoyer@neon-lab.com>"

start on runlevel [234]
stop on runlevel [0156]

kill timeout 600
limit nofile 65536 65536

pre-start script
    [ -d /var/run/flume-ng/flume-ng-agent ] || mkdir -p /var/run/flume-ng/flume-ng-agent 
    chown -R <%= @flume_user %>:<%= @flume_user %> /var/run/flume-ng/    
    chdir /var/run/flume-ng/flume-ng-agent
end script

script
    cd /var/run/flume-ng/flume-ng-agent
    exec sudo -u <%= @flume_user %> <%= @flume_bin %> agent --conf <%= @flume_conf_dir %> --conf-file <%= @flume_conf_dir %>/flume.conf --name agent
end script
respawn
