# video_client/templates/default/video_client_service.conf.erb
description "Neon Video Client"

start on runlevel [234]
stop on runlevel [0156]

kill timeout 1800
limit nofile 65536 65536

setuid <%= @user %>
setgid <%= @group %>

chdir <%= @neon_root_dir %>
env VIRTUAL_ENV=<%= @neon_root_dir %>/.pyenv
env PATH="${VIRTUAL_ENV}/bin:$PATH"
exec ${VIRTUAL_ENV}/bin/python <%= @neon_root_dir %>/api/client.py -c <%= @config_file%>

respawn
