# stats_manager/templates/default/statsmanager_service.conf.erb
description "Neon Stats Manager"
author  "Mark Desnoyer <desnoyer@neon-lab.com>"

start on runlevel [234]
stop on runlevel [0156]

kill timeout 30
limit nofile 65536 65536

setuid <%= @user %>
setgid <%= @group %>

chdir <%= @neon_root_dir %>
exec <%= @neon_root_dir %>/stats/cluster_manager.py -c <%= @config_file%> --stats.batch_processor.mr_jar <%= @mr_jar %>
respawn