# neon/templates/default/trackserver_service.conf.erb
description "Neon Click Trackserver"
author  "Mark Desnoyer <desnoyer@neon-lab.com>"

start on runlevel [234]
stop on runlevel [0156]

kill timeout 30

setuid <%= user %>
setgid <%= group %>

chdir <%= neon_root_dir %>
exec <%= neon_root_dir %>/clickTracker/trackserver.py -c <%= config_file%>
respawn